# crontab
@reboot python3 /home/warproxxx/IPFS/runall.py > logs.txt